# systemd/power-monitor.service
[Service]
Environment=STARTUP_WAIT=30s
ExecStart=%h/.config/systemd/user/power_monitor.sh

[Install]
WantedBy=default.target
